# my_twitter.conf
description "Gunicorn application server running my_twitter"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env PATH=/home/ec2-user/my_twitter/env/bin
env PROGRAM_NAME="my_twitter"
env USERNAME="ec2-user"
env GUNICORN_BIN=/home/ec2-user/my_twitter/env/bin/gunicorn
env SOCKET_PATH=/var/run/my-twitter.sock

# Main script to be run
script
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Ready to run..." >> /var/log/$PROGRAM_NAME.sys.log

    export HOME="/home/ec2-user"
    export REDIS_URL="localhost"
    export REDIS_PORT=6379
    export SECRET_KEY='dev'
    export OAUTHLIB_INSECURE_TRANSPORT='1'
    echo $$ > /var/run/$PROGRAM_NAME.pid

    cd /home/ec2-user/my_twitter/src/my_twitter
    # exec sudo -u ec2-user $GUNICORN_BIN --workers 3 --bind unix:my_twitter.sock -m 000 wsgi >> /var/log/$PROGRAM_NAME.sys.log 2>&1
    # exec su -s /bin/sh -c 'exec "$0" "$@"' ec2-user -- $GUNICORN_BIN --workers 3 --bind unix:my_twitter.sock -m 000 wsgi >> /var/log/$PROGRAM_NAME.sys.log 2>&1
    exec $GUNICORN_BIN --workers 3 --bind unix:$SOCKET_PATH -m 000 wsgi >> /var/log/$PROGRAM_NAME.sys.log 2>&1
end script

# Script for debug purpose, run before starting
pre-start script
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/$PROGRAM_NAME.sys.log
end script

# Script for debug purpose, run before stopping
pre-stop script
    rm /var/run/$PROGRAM_NAME.pid/
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/$PROGRAM_NAME.sys.log
end script